# etc/udev/98-hdmi-hotplug.rules
SUBSYSTEM=="drm", ACTION=="change", \
    ENV{DISPLAY}=":0", ENV{XAUTHORITY}="/home/work/.Xauthority", RUN+="/bin/sh /root/.bin/rw-hdmi-hotplug.sh"
