# etc/init.d/v2ray
#!/sbin/openrc-run
# Copyright 1999-2015 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

V2RAY_CONFIG="/etc/v2ray.json"
V2RAY_COMMAND="/usr/local/bin/v2ray"
pidfile="/run/${RC_SVCNAME}.pid"

depend() {
    need net
}

checkconfig() {
    if [ ! -f ${V2RAY_CONFIG} ]; then
        ewarn "${V2RAY_CONFIG} does not exist."
        return 1
    fi
}

start() {
    checkconfig || return 1

    ebegin "Starting V2ray"
    start-stop-daemon --start -m --pidfile $pidfile --exec ${V2RAY_COMMAND} \
    -- -config ${V2RAY_CONFIG} >/dev/null 2>&1 &
    eend $?
}

stop() {
    ebegin "Stopping V2ray"
    start-stop-daemon --stop --pidfile ${pidfile}
    eend $?
}

